# etc/cron/crontab
* * * * * sh /variables.sh && /tutorial/tutorial -app=script -name=test
